# .docker/aiida-core-with-services/s6-assets/s6-rc.d/postgresql/up
#!/command/execlineb -P

with-contenv

pg_ctl -D /home/aiida/.postgresql -l /home/${SYSTEM_USER}/.postgresql/logfile start
