# .docker/aiida-core-base/s6-assets/s6-rc.d/run-before-daemon-start/up
#!/command/execlineb -P

with-contenv

foreground { s6-echo "Calling /etc/init/run-before-daemon-start" }
run-parts --regex=".*" /etc/init/run-before-daemon-start/
