# .docker/aiida-core-with-services/s6-assets/s6-rc.d/rabbitmq-init/up
#!/command/execlineb -S0
with-contenv

foreground { s6-echo "Calling /etc/init/rabbitmq-init.sh" }
/etc/init/rabbitmq-init.sh
