# .docker/aiida-core-base/s6-assets/s6-rc.d/aiida-daemon-start/up
#!/command/execlineb -S0

verdi daemon start
