# .docker/aiida-core-base/s6-assets/s6-rc.d/aiida-daemon-start/up
#!/command/execlineb -S0

with-contenv

foreground { s6-echo "Calling /etc/init/aiida-daemon-start" }
/etc/init/aiida-daemon-start.sh
