# .docker/aiida-core-with-services/s6-assets/s6-rc.d/rabbitmq/run
#!/command/execlineb -P

with-contenv

foreground { s6-echo "Calling /etc/init/rabbitmq.sh" }
rabbitmq-server
