# .docker/base-with-services/s6-assets/s6-rc.d/postgresql/run
#!/command/execlineb -P

with-contenv

# # -w waits until server is up
# echo "DEBUG"
# echo "${SYSTEM_USER}"
#
# PSQL_START_CMD="pg_ctl --timeout=180 -w -D /home/${SYSTEM_USER}/.postgresql -l /home/${SYSTEM_USER}/.postgresql/logfile start"
# PSQL_STATUS_CMD="pg_ctl -D /home/${SYSTEM_USER}/.postgresql status"
#
# MAMBA_RUN="mamba run -n aiida-core-services"
#
# # Fix problem with kubernetes cluster that adds rws permissions to the group
# # for more details see: https://github.com/materialscloud-org/aiidalab-z2jh-eosc/issues/5
# chmod g-rwxs /home/${SYSTEM_USER}/.postgresql -R
#
# # stores return value in $?
# running=true
# ${MAMBA_RUN} ${PSQL_STATUS_CMD} || running=false
#
# # Postgresql was probably not shutdown properly. Cleaning up the mess...
# if ! $running ; then
#    echo "" > /home/${SYSTEM_USER}/.postgresql/logfile # empty log files
#    rm -vf /home/${SYSTEM_USER}/.postgresql/postmaster.pid
#    ${MAMBA_RUN} ${PSQL_START_CMD}
# fi

mamba run -n aiida-core-services postgres -D /home/aiida/.postgresql
