# .docker/aiida-core-base/s6-assets/s6-rc.d/run-after-daemon-start/up
#!/command/execlineb -P

with-contenv

foreground { s6-echo "Calling /etc/init/run-after-daemon-start" }
run-parts --regex=".*" /etc/init/run-after-daemon-start/
